# templates/node_exporter.service.j2
[Unit]
Description=prometheus node exporter
Requires=network-online.target
After=network-online.target

[Service]
PrivateTmp={{ node_exporter_private_tmp }}
User={{ node_exporter_user }}
Group={{ node_exporter_group }}
RuntimeDirectory=node_exporter
EnvironmentFile={{ node_exporter_environment_path }}/node_exporter
ExecStart={{ node_exporter_bin_path }}/node_exporter {% if _node_exporter_collectors_as_string_version %} {{ _node_exporter_config_param_selector }}collectors.enabled $NODE_EXPORTER_COLLECTORS_ENABLED {% else %} {% for option in node_exporter_collectors_enabled.split(',') %}{{ _node_exporter_config_param_selector }}collector.{{option}} {% endfor %} {% endif %} {{ _node_exporter_config_param_selector }}collector.textfile.directory $NODE_EXPORTER_COLLECTOR_TEXTFILE_DIRECTORY {{ _node_exporter_config_param_selector }}collector.diskstats.ignored-devices $NODE_EXPORTER_COLLECTOR_DISKSTATS_IGNORED {{ _node_exporter_config_param_selector }}collector.filesystem.ignored-fs-types $NODE_EXPORTER_COLLECTOR_FILESYSTEM_IGNORED {{ _node_exporter_config_param_selector }}web.listen-address :{{ node_exporter_port }} {{ _node_exporter_config_param_selector }}web.telemetry-path /{{ node_exporter_path }} {{ _node_exporter_config_param_selector }}log.level $NODE_EXPORTER_LOG_LEVEL {% for option in node_exporter_config_extra_options %}{{ _node_exporter_config_param_selector }}{{ option }} {% endfor %}

StandardOutput={{ node_exporter_log_output }}
StandardError={{ node_exporter_log_output }}
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM
Restart=always

[Install]
WantedBy=multi-user.target
